# etc/cleanblog.conf
[program:cleanblog]
command = /home/xin/workspace/python/cleanblog/venv/bin/gunicorn -b 127.0.0.1:9000 wsgi_gunicorn:app
directory = /home/xin/workspace/python/cleanblog
user = shin
